# crontab
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
* * * * * /bin/bash /app/wrapper.sh >> /var/log/cron.log 2>&1

